# playbook/oprspedia.service.j2
[Unit]
Description=Oprspedia

[Service]
PIDfile={{ settings.pidfile }}
ExecStart=thin -C {{ source.path }}/src/oprspedia.ru -C {{ source.path }}/src/config.yaml start

[Install]
WantedBy=multi-user.target
