# playbook/oprspedia.service.j2
[Unit]
Description=Oprspedia

[Service]
ExecStart=/usr/bin/ruby {{ source.path }}/src/oprspedia.rb

[Install]
WantedBy=multi-user.target
